# rootfs/etc/s6-overlay/s6-rc.d/gost/run
#!/command/with-contenv sh
set -e
exec /usr/local/bin/gost -D -C /config/gost.yaml
