# rootfs/etc/s6-overlay/s6-rc.d/easytier/run
#!/command/with-contenv sh
set -e
exec /usr/local/bin/easytier-core -c /config/config.toml
